# systemd/duplicity-unattended.timer
[Unit]
Description=Run daily Duplicity backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=20min

[Install]
WantedBy=timers.target
