# systemd/duplicity-unattended.service
[Unit]
Description=Unattended Duplicity backup

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/duplicity-unattended
